# contrib/docker/ckan_harvesting.conf
; ===============================
; ckan harvester
; ===============================

[program:ckan_gather_consumer]

command=/usr/local/bin/ckan-paster --plugin=ckanext-harvest harvester gather_consumer --config=/etc/ckan/production.ini

; user that owns virtual environment.
user=root

numprocs=1
stdout_logfile=/var/log/ckan/std/gather_consumer.log
stderr_logfile=/var/log/ckan/std/gather_consumer.log
autostart=true
autorestart=true
startsecs=10

[program:ckan_fetch_consumer]

command=/usr/local/bin/ckan-paster --plugin=ckanext-harvest harvester fetch_consumer --config=/etc/ckan/production.ini

; user that owns virtual environment.
user=root

numprocs=1
stdout_logfile=/var/log/ckan/std/fetch_consumer.log
stderr_logfile=/var/log/ckan/std/fetch_consumer.log
autostart=true
autorestart=true
startsecs=10
